# deploy/drillup.service
[Unit]
Description=drillup Next.js app
After=network.target

[Service]
Type=simple
WorkingDirectory=__DEPLOY_PATH__
ExecStart=__NPM_BIN__ start
Restart=on-failure
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
